# systemd/hwt-cleanup.service
[Unit]
Description=HWT905 Data Cleanup Service
Documentation=https://github.com/your-repo/hwt905-raspi-lite

[Service]
Type=oneshot

# Đường dẫn tới file môi trường
EnvironmentFile=/home/pi/hwt905-raspi-lite/.env

# Người dùng để chạy service
User=pi
Group=pi

# Thư mục làm việc
WorkingDirectory=/home/pi/hwt905-raspi-lite

# Lệnh để chạy kịch bản cleanup
ExecStart=/home/pi/hwt905-raspi-lite/venv/bin/python3 scripts/cleanup.py

# Ghi log vào journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hwt-cleanup 